# com.vtmeilab.paperbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS launchd 定时任务配置。
  安装方法见 README「方式 B：launchd」。
  下面写死了每天 08:30 运行，改 Hour/Minute 即可。
  注意：把 __DIR__ 替换成 run.sh 的真实绝对路径。
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vtmeilab.paperbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__DIR__/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__DIR__/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__DIR__/launchd.err.log</string>
</dict>
</plist>
